# planner/rootfs/etc/services.d/planner/run
#!/usr/bin/with-contenv bash
exec bash /app/run.sh
