# planner_autonome/rootfs/etc/services.d/planner/finish
#!/usr/bin/with-contenv bash
echo "[PLANNER] Service stopped"
exit 0
