# planner_autonome/rootfs/etc/services.d/planner/finish
#!/bin/bash
echo "[Planner] Service stopped"
