# planner_autonome/rootfs/etc/services.d/planner/run
#!/usr/bin/with-contenv bash
cd /opt/planner/service
echo "[PLANNER] Starting server..."
exec node server.js
