# planner_autonome/rootfs/etc/services.d/planner/run
#!/usr/bin/execlineb -P
with-contenv
cd /opt/planner/service
node server.js
